# shoutcast/Supervisord/supervisord.conf
[supervisord]
nodaemon=false
logfile = /var/log/messages/supervisord.log

[program:tcpdump]
command=/bin/bash /usr/local/shoutcast/tcpdump.sh
autorestart=true